# scenario1/root/etc/containers/systemd/odoo-db.container
[Unit]
Description=PostgreSQL Database Server
Documentation=https://www.postgresql.org/
After=network.target

[Container]
ContainerName=odoo-db
Image=docker.io/library/postgres:17-alpine

# Network configuration
Network=host

# Security context (equivalent to K8s securityContext)
#NoNewPrivileges=true
#DropCapability=ALL

# Environment variables from config
EnvironmentFile=/etc/containers/systemd/configs/odoo-db.env

# Volume mounts
Volume=/var/lib/postgresql/data:/var/lib/postgresql/data:Z

# Health check
HealthCmd=pg_isready -U odoo
HealthInterval=30s
HealthTimeout=10s
HealthStartPeriod=60s
HealthRetries=3

[Service]
Restart=always
RestartSec=10
TimeoutStartSec=120
TimeoutStopSec=30

[Install]
WantedBy=odoo.target
